# app/app.conf
description "My application server"
author "root@localhost"
start on runlevel [2345]
stop on runlevel [!2345]
env AN_ENVIRONMENTAL_VARIABLE=i-want-to-set
respawn
exec /vagrant/app/app.py

